# supervisor/config.d/gunicorn.ini
[program:gunicorn]
user=root
directory=/root/www/warehouse
command=/root/.virtualenvs/leo/bin/gunicorn --pythonpath /root/.virtualenvs/leo/bin/ -c /root/www/warehouse/gunicorn/gunicorn-conf.py warehouse.wsgi
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/root/www/warehouse/gunicorn/out.log
stderr_logfile=/root/www/warehouse/gunicorn/error.log

