# supervisor/config.d/nginx.ini
[program:nginx]
user=root
command=/usr/local/nginx/sbin/./nginx -c /root/www/warehouse/nginx/nginx.conf
autostart=true
autorestart=true
startretries=3
exitcodes=0
stopwaitsecs=0
redirect_stderr=true
